# systemd/usb-importer.service
[Unit]
Description=USB WAV importer web/api
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/voicenotes
Environment=CONFIG=/var/lib/voicenotes/config.yaml
ExecStart=/opt/voicenotes/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8080
Restart=on-failure

[Install]
WantedBy=multi-user.target